# setup/systemd/cpt-ai-api.service
[Unit]
Description=CPT AI Regression Analyzer - REST API
After=network.target cpt-ai-mcp.service

[Service]
Type=simple
User=__USER__
WorkingDirectory=__REPO_DIR__/components/api
ExecStart=__VENV_DIR__/bin/python3 -m uvicorn api:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5
EnvironmentFile=-__REPO_DIR__/.env

[Install]
WantedBy=multi-user.target
