# setup/systemd/cpt-ai-mcp.service
[Unit]
Description=CPT AI - OpenSearch MCP Server
After=network.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__REPO_DIR__
ExecStart=__VENV_DIR__/bin/opensearch-mcp-server-py --mode multi --config __HOME__/opensearch_mcp_config.yml --transport stream --port 9900
Restart=on-failure
RestartSec=5
EnvironmentFile=-__REPO_DIR__/.env

[Install]
WantedBy=multi-user.target
