# setup/systemd/cpt-ai-dashboard.service
[Unit]
Description=CPT AI - RHEL Performance Engineering Dashboard
After=network.target cpt-ai-mcp.service

[Service]
Type=simple
User=__USER__
WorkingDirectory=__REPO_DIR__/components/dashboard
ExecStart=__VENV_DIR__/bin/python3 app.py
Restart=on-failure
RestartSec=5
Environment=PORT=8050
EnvironmentFile=-__REPO_DIR__/.env

[Install]
WantedBy=multi-user.target
